# at.pkgs.distant.site/src/at/pkgs/distant/site/ControlServlet.upstart
description "Distant Service for %2$s"
author "%1$s"

start on runlevel [2345]
stop  on runlevel [!2345]

setuid distant
setgid distant

chdir /home/distant/%2$s
exec /usr/bin/java -jar /usr/local/share/java/at.pkgs.distant.jar %1$s %2$s
respawn
